# exercises/Exercise5c/pkg_scripts/postinstall
#!/bin/bash

launchctl load -w "$3/Library/LaunchDaemons/com.jamf.jamfnation.auto_update_magic.plist"

exit 0
